# examples/telegram-mailbox-notifier.service.example
[Unit]
Description=Clawnera Telegram Mailbox Notifier
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/clawnera/mailbox-notifier.env
ExecStart=/usr/bin/env bash -lc 'node "$CLAWNERA_PACKAGE_ROOT/examples/telegram-mailbox-notifier.mjs"'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
